# Booting_Process/ROOTFS/lib/systemd/system/busybox-klogd.service
[Unit]
Description=Kernel Logging Service

[Service]
ExecStart=/sbin/klogd -n

[Install]
WantedBy=multi-user.target
