# 01_Booting_Process/ROOTFS/lib/systemd/system/connman.service
[Unit]
Description=Connection service
After=syslog.target

[Service]
Type=dbus
BusName=net.connman
ExecStart=/usr/sbin/connmand -n

[Install]
WantedBy=multi-user.target
