# Booting_Process/ROOTFS/var/lib/opkg/info/avahi-daemon.postinst
if [ -z "$D" ]; then
	DBUSPID=`pidof dbus-daemon`

	if [ "x$DBUSPID" != "x" ]; then
		/etc/init.d/dbus-1 force-reload
	fi
fi
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT="-s"
fi
update-rc.d $OPT avahi-daemon defaults 21 19
