# Booting_Process/ROOTFS/var/lib/opkg/info/avahi-systemd.postrm
#!/bin/sh
systemctl disable avahi-daemon.service
