# Booting_Process/ROOTFS/var/lib/opkg/info/connman.postrm
#!/bin/sh
update-rc.d $D connman remove
