# Booting_Process/ROOTFS/var/lib/opkg/info/connman.list
/etc/dbus-1/system.d/connman.conf
/etc/init.d/connman
/usr/sbin/connmand
