# Booting_Process/ROOTFS/var/lib/opkg/info/dropbear.list
/usr/bin/dbclient
/usr/sbin/dropbearmulti
/etc/dropbear/
/etc/init.d/dropbear
/usr/sbin/dropbearkey
/usr/sbin/dropbearconvert
/usr/sbin/dropbear
/etc/pam.d/dropbear
