# Booting_Process/ROOTFS/var/lib/opkg/info/dropbear-systemd.postrm
#!/bin/sh
systemctl disable dropbear.socket
