# Booting_Process/ROOTFS/var/lib/opkg/info/netbase.postrm
#!/bin/sh
update-rc.d $D networking remove
