# Booting_Process/ROOTFS/var/lib/opkg/info/systemd-compat-units.list
/lib/systemd/system/basic.target.wants/run-postinsts.service
/lib/systemd/system/basic.target.wants/machineid.service
/lib/systemd/system/run-postinsts.service
/lib/systemd/system/sysinit.target.wants/machineid.service
/lib/systemd/system/machineid.service
/lib/systemd/system/sysinit.target.wants/run-postinsts.service
/usr/bin/runlevel
