# Booting_Process/ROOTFS/var/lib/opkg/info/wpa-supplicant-systemd.postrm
#!/bin/sh
systemctl disable wpa_supplicant.service
