# 12_I2C_RTC_DS3231/ebb_clock.service
[Unit]
Description=EBB RTC Service

[Service]
Type=simple
ExecStart=/usr/bin/rtcds3231 1

[Install]
WantedBy=multi-user.target