# 01_Booting_Process/ROOTFS/var/lib/opkg/info/avahi-daemon.postrm
	deluser avahi || true
	delgroup avahi || true
update-rc.d $D avahi-daemon remove
