# 01_Booting_Process/ROOTFS/var/lib/opkg/info/avahi-systemd.list
/lib/systemd/system/avahi-daemon.service
/lib/systemd/system/avahi-dnsconfd.service
/lib/systemd/system/avahi-daemon.socket
