# 01_Booting_Process/ROOTFS/var/lib/opkg/info/avahi-systemd.postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

systemctl $OPTS enable avahi-daemon.service

if [ -z "$D" -a enable = "enable" ]; then
    systemctl start avahi-daemon.service
fi
