# 01_Booting_Process/ROOTFS/var/lib/opkg/info/bluez4.list
/etc/bluetooth/network.conf
/etc/bluetooth/input.conf
/usr/lib/libbluetooth.so.3.13.0
/usr/sbin/bluetoothd
/usr/share/dbus-1/system-services/org.bluez.service
/lib/udev/hid2hci
/lib/systemd/system/bluetooth.service
/usr/bin/ciptool
/usr/sbin/hciattach
/usr/bin/l2ping
/usr/bin/hcitool
/etc/bluetooth/rfcomm.conf
/usr/bin/sdptool
/etc/bluetooth/main.conf
/usr/bin/dfutool
/usr/bin/l2test
/usr/lib/libbluetooth.so.3
/etc/dbus-1/system.d/bluetooth.conf
/usr/bin/rctest
/etc/bluetooth/
/lib/udev/rules.d/97-bluetooth-hid2hci.rules
/etc/bluetooth/audio.conf
/usr/sbin/hciemu
/usr/sbin/hciconfig
/usr/bin/gatttool
/usr/bin/rfcomm
/usr/sbin/bccmd
