# 01_Booting_Process/ROOTFS/var/lib/opkg/info/busybox-syslog-systemd.postrm
#!/bin/sh
systemctl disable busybox-syslog.service
