# 01_Booting_Process/ROOTFS/var/lib/opkg/info/busybox-syslog.list
/etc/syslog-startup.conf.busybox
/etc/init.d/syslog.busybox
