# 01_Booting_Process/ROOTFS/var/lib/opkg/info/busybox-syslog.postrm
#!/bin/sh
	update-alternatives --remove  syslog-init syslog.busybox
	update-alternatives --remove  syslog-startup-conf syslog-startup.conf.busybox
update-rc.d $D syslog remove
