# 01_Booting_Process/ROOTFS/var/lib/opkg/info/connman-systemd.postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

systemctl $OPTS enable connman.service

if [ -z "$D" -a enable = "enable" ]; then
    systemctl start connman.service
fi
