# 01_Booting_Process/ROOTFS/var/lib/opkg/info/connman-systemd.postrm
#!/bin/sh
systemctl disable connman.service
