# 01_Booting_Process/ROOTFS/var/lib/opkg/info/dbus-1.postinst
	if [ -z "$D" ] && [ -e /etc/init.d/populate-volatile.sh ] ; then
		/etc/init.d/populate-volatile.sh update
	fi
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT="-s"
fi
update-rc.d $OPT dbus-1 start 02 5 3 2 . stop 20 0 1 6 .
