# 01_Booting_Process/ROOTFS/var/lib/opkg/info/dropbear.postrm
#!/bin/sh
	update-alternatives --remove  scp ../sbin/dropbearmulti
	update-alternatives --remove  ssh ../sbin/dropbearmulti
update-rc.d $D dropbear remove
