# 01_Booting_Process/ROOTFS/var/lib/opkg/info/ofono.list
/lib/systemd/system/ofono.service
/etc/ofono/
/etc/dbus-1/system.d/ofono.conf
/usr/sbin/ofonod
/etc/ofono/phonesim.conf
/etc/init.d/ofono
