# 01_Booting_Process/ROOTFS/var/lib/opkg/info/systemd-serialgetty.list
/etc/systemd/system/getty.target.wants/serial-getty@ttyO2.service
/lib/systemd/system/serial-getty@.service
