# 01_Booting_Process/ROOTFS/var/lib/opkg/info/wpa-supplicant-systemd.postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

systemctl $OPTS disable wpa_supplicant.service

if [ -z "$D" -a disable = "enable" ]; then
    systemctl start wpa_supplicant.service
fi
