# 01_Booting_Process/ROOTFS/var/lib/opkg/info/wpa-supplicant-systemd.list
/lib/systemd/system/wpa_supplicant-nl80211@.service
/lib/systemd/system/wpa_supplicant-wired@.service
/lib/systemd/system/wpa_supplicant@.service
/lib/systemd/system/wpa_supplicant.service
